# scripts/launchd/com.betvector.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  BetVector — Morning Pipeline (PC-15-04)
  ========================================
  Runs at 07:00 local time daily.
  Full pipeline: scrape fixtures + odds → compute features → generate
  predictions → detect value bets → send morning email.

  Install:  cp this file ~/Library/LaunchAgents/ && launchctl load ~/Library/LaunchAgents/com.betvector.morning.plist
  Unload:   launchctl unload ~/Library/LaunchAgents/com.betvector.morning.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.betvector.morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kyng/Documents/BetVector/scripts/run_pipeline_local.sh</string>
        <string>morning</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kyng/Documents/BetVector</string>

    <!-- Run at 07:00 local time every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Log stdout and stderr to data/logs/ -->
    <key>StandardOutPath</key>
    <string>/Users/kyng/Documents/BetVector/data/logs/launchd_morning_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kyng/Documents/BetVector/data/logs/launchd_morning_stderr.log</string>

</dict>
</plist>
