# scripts/launchd/com.betvector.evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  BetVector — Evening Pipeline (PC-15-04)
  ========================================
  Runs at 21:00 local time daily.
  Resolves settled bets against actual match results, updates bankroll
  P&L, and sends the evening summary email.

  Install:  cp this file ~/Library/LaunchAgents/ && launchctl load ~/Library/LaunchAgents/com.betvector.evening.plist
  Unload:   launchctl unload ~/Library/LaunchAgents/com.betvector.evening.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.betvector.evening</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kyng/Projects/BetVector/scripts/run_pipeline_local.sh</string>
        <string>evening</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kyng/Projects/BetVector</string>

    <!-- Run at 21:00 local time every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Log stdout and stderr to data/logs/ -->
    <key>StandardOutPath</key>
    <string>/Users/kyng/Projects/BetVector/data/logs/launchd_evening_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kyng/Projects/BetVector/data/logs/launchd_evening_stderr.log</string>

</dict>
</plist>
