# scripts/com.betvector.wc_evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  BetVector — WC Evening Pipeline (WC-07-02)  [RETIRED — WC-10-02]
  ================================================================
  *** DO NOT INSTALL. *** Retired 2026-06-23: 22:00 ET fired while the 10/11 PM ET
  matches were still being played, so it mis-captured results. Its duties (results,
  CLV capture, model accuracy) are now absorbed by the morning run (09:30 ET), which
  is the single daily spine. The `evening` pipeline mode itself still works for a
  manual "settle results now" run; only the *schedule* is retired.

  Manual run (optional):  scripts/run_wc_pipeline.sh evening
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.betvector.wc_evening</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kyng/Projects/BetVector/scripts/run_wc_pipeline.sh</string>
        <string>evening</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kyng/Projects/BetVector</string>

    <!-- Run at 22:00 local time every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/kyng/Projects/BetVector/data/logs/launchd_wc_evening_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kyng/Projects/BetVector/data/logs/launchd_wc_evening_stderr.log</string>

</dict>
</plist>
