# scripts/com.betvector.wc_evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  BetVector — WC Evening Results Refresh (WC-07-02; re-enabled + retimed 2026-06-25)
  =================================================================================
  A LIGHT same-day results settle so evening matches show up that night instead of
  waiting for the 09:30 ET morning run. Runs `run_wc_pipeline.sh evening`, i.e. the
  pipeline's evening mode: scrape results -> capture closing-line value (CLV) ->
  update Elo -> compute model accuracy -> evening email. It does NOT re-pull odds,
  re-predict, or scan for value bets (that's the morning run's job) — so it's light.

  RETIMED to 23:30 ET (was 22:00 ET): the old 22:00 slot fired while the late 10/11
  PM ET matches were still in play and mis-captured them. 23:30 ET runs after the
  day's matches have finished; the 09:30 ET morning run still backstops any rare late
  finisher (it pulls the last 3 days of scores).

  Install:  cp scripts/com.betvector.wc_evening.plist ~/Library/LaunchAgents/ && launchctl load ~/Library/LaunchAgents/com.betvector.wc_evening.plist
  Unload:   launchctl unload ~/Library/LaunchAgents/com.betvector.wc_evening.plist
  Manual:   scripts/run_wc_pipeline.sh evening
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.betvector.wc_evening</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kyng/Projects/BetVector/scripts/run_wc_pipeline.sh</string>
        <string>evening</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kyng/Projects/BetVector</string>

    <!-- Run at 23:30 ET every day — after the day's matches have finished -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/kyng/Projects/BetVector/data/logs/launchd_wc_evening_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kyng/Projects/BetVector/data/logs/launchd_wc_evening_stderr.log</string>

</dict>
</plist>
